# com.yuma.property-patrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuma.property-patrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/.venv/bin/python</string>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/run_daily_patrol.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ytejima/Documents/Projects/property-analyzer</string>

    <!-- Fixed clock: fires at 01:00 AM every night regardless of when last run completed.
         StartInterval (relative 24h) caused drift: if a run completed at 9AM, the next
         fire was 9AM next day → idle guard spun 2h → completion at 11AM → cascade worse.
         StartCalendarInterval + night-hours idle_guard = always runs clean at 1:00 AM. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- RunAtLoad disabled: was triggering extra runs on reboot/login during daytime,
         which started the drift cascade. Manual catch-up: run_daily_patrol.py --force -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/property_patrol.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/property_patrol.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/ytejima/Documents/Projects</string>
    </dict>
</dict>
</plist>
