# com.yuma.property-patrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuma.property-patrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/.venv/bin/python</string>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/run_daily_patrol.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ytejima/Documents/Projects/property-analyzer</string>

    <!-- Fires at 17:00, not 01:00. The Mac is asleep at night (lid closed / carried /
         power unplugged, ~18:00 through early morning) — 01:00 was structurally
         undeliverable, not just occasionally flaky (see project_laptop_evening_unavailable
         memory; same root cause already fixed once for daily-digest's 18:00→16:00 move).
         17:00 sits inside the confirmed-awake 16:00-17:30 window with runway before
         ~18:00 sleep onset (typical run ≈30min per data/patrol_summary.json elapsed_min).
         GHA daily-patrol.yml's 06:00 JST cron schedule: is the real fallback for any
         day this still doesn't fire — see that workflow's own comment for the history. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- RunAtLoad disabled: was triggering extra runs on reboot/login during daytime,
         which started the drift cascade. Manual catch-up: run_daily_patrol.py --force -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/property_patrol.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/property_patrol.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/ytejima/Documents/Projects</string>
    </dict>
</dict>
</plist>
