# com.yuma.property-patrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuma.property-patrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/.venv/bin/python</string>
        <string>/Users/ytejima/Documents/Projects/property-analyzer/run_daily_patrol.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ytejima/Documents/Projects/property-analyzer</string>

    <!-- DISABLED as of 2026-07-08 (not installed in ~/Library/LaunchAgents). GHA
         daily-patrol.yml is now Primary (06:00 JST cron, unconditional, repo is public
         so Actions minutes are free). Local used to be Primary but the Mac is asleep
         most hours (see project_laptop_evening_unavailable memory) — even the best
         reachable local window (17:00) is still probabilistic (Yuma might be out).
         Cloud has no such failure mode. This plist is kept only as a reference for
         manual/opportunistic runs: python run_daily_patrol.py --force. No automatic
         trigger key is set on purpose — copying this file back into LaunchAgents does
         nothing unless a StartCalendarInterval/RunAtLoad key is re-added deliberately. -->

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/property_patrol.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/property_patrol.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/ytejima/Documents/Projects</string>
    </dict>
</dict>
</plist>
