# scripts/com.yuma.property-daily-email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuma.property-daily-email</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yumatejima/Documents/Projects/stock-analyzer/.venv/bin/python3</string>
        <string>/Users/yumatejima/Documents/Projects/property-analyzer/scripts/send_daily_property_email.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/yumatejima/Documents/Projects</string>
    <key>StandardOutPath</key>
    <string>/tmp/property_daily_email.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/property_daily_email.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
